# linkvision.service
[Unit]
Description=LinkVision - Network Infrastructure Visualization
After=network.target

[Service]
# Запуск от root (требуется для raw ICMP и сетевых операций)
User=root
Group=root

# Рабочая директория проекта
WorkingDirectory=/opt/linkvision

# Команда запуска через Gunicorn
# -w 4 : количество воркеров (процессов)
# -b 0.0.0.0:5000 : адрес и порт прослушивания
# app:app : файл app.py и объект приложения app внутри него
ExecStart=/opt/linkvision/venv/bin/gunicorn -w 4 -b 0.0.0.0:5000 app:app

# Переменные окружения
Environment="PATH=/opt/linkvision/venv/bin"
Environment="FLASK_ENV=production"
Environment="FLASK_APP=app.py"
# Если используете .env-файл:
# EnvironmentFile=/opt/linkvision/.env

# Автоматический перезапуск при падении
Restart=always
RestartSec=10

# Логирование в systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=linkvision

# Опционально: базовые ограничения безопасности
PrivateTmp=true

[Install]
WantedBy=multi-user.target